# rootfs/etc/s6-overlay/s6-rc.d/init-thermostart/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Thermostart
# Preparing configuration for Thermostart
# ==============================================================================
export DATABASE_PATH=/config/thermostart/thermostart.db
export DATABASE_URL=sqlite:////config/thermostart/thermostart.db
export APP_FOLDER=/opt/services/web
export FLASK_APP=thermostart/__init__.py

if bashio::config.has_value 'secret_key'; then
    export SECRET_KEY="$(bashio::config 'secret_key')"
fi

if ! bashio::fs.file_exists "$DATABASE_PATH"; then
      mkdir -p /config/thermostart
fi

cd /opt/services/web

python3 manage.py has_alembic_version
if (test $? -ne 1)
then
    python3 manage.py db stamp 26422f1f63d0
fi
python3 manage.py db upgrade
python3 manage.py fill_db
exit 0