# rootfs/etc/s6-overlay/s6-rc.d/thermostart/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Thermostart
# Runs Thermostart
# ==============================================================================
bashio::log.info "Starting Thermostart..."

export FLASK_APP=thermostart/__init__.py
export DATABASE_URL=sqlite:////config/thermostart/thermostart.db
export APP_FOLDER=/opt/services/web
export APPLICATION_ROOT=$(bashio::addon.ingress_entry)
export SCRIPT_NAME=$(bashio::addon.ingress_entry)

if bashio::config.has_value 'secret_key'; then
    export SECRET_KEY="$(bashio::config 'secret_key')"
fi

if bashio::config.has_value 'autologin.username'; then
    export AUTOLOGIN_USERNAME="$(bashio::config 'autologin.username')"
fi

if bashio::config.has_value 'autologin.password'; then
    export AUTOLOGIN_PASSWORD="$(bashio::config 'autologin.password')"
fi

cd /opt/services/web

exec \
  gunicorn -k eventlet \
    -w 1 \
    -b 0.0.0.0:3889 \
    'thermostart:create_app()'
